# packages/service/tvheadend/system.d/tvheadend.service
[Unit]
Description=tvheadend
ConditionPathExists=|/storage/.cache/services/tvheadend.conf
After=multi-user.target

[Service]
EnvironmentFile=-/storage/.cache/services/tvheadend.conf
ExecStart=/bin/sh -c 'exec /usr/bin/tvheadend.start'
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target