# packages/addons/service/vdr-addon/source/system.d/service.multimedia.vdra.service
[Unit]
Description=vdr
ConditionPathExists=|/storage/.kodi/userdata/addon_data/service.multimedia.vdra/vdr.enable
ConditionPathExists=|/storage/.kodi/userdata/addon_data/service.multimedia.vdra/vdr-front.run
After=graphical.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.multimedia.vdra/bin/vdra.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
