# packages/addons/service/udocker/source/system.d/service.system.udocker.socket
[Unit]
Description=Docker Socket for the API
PartOf=/storage/.kodi/addons/service.system.udocker/system.d/service.system.udocker.service

[Socket]
ListenStream=/var/run/docker.sock
SocketMode=0660
SocketUser=root
SocketGroup=root

[Install]
WantedBy=sockets.target