# packages/services/oscam-service/system.d/oscam.service
[Unit]
Description=oscam
ConditionPathExists=|/storage/.cache/services/oscam.conf
After=multi-user.target

[Service]
ExecStart=/bin/sh -c "exec sh /usr/bin/oscam.start"
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target