# packages/addons/service/vdr24/source/system.d/service.vdr24.service
[Unit]
Description=vdr
ConditionPathExists=|/storage/.kodi/userdata/addon_data/service.vdr24/vdr.enable
ConditionPathExists=|/storage/.kodi/userdata/addon_data/service.vdr24/vdr-front.run
After=graphical.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.vdr24/bin/vdr.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
