# packages/addons/service/embyu/source/system.d/service.embyu.service
[Unit]
Description=Emby - streaming server (Ubuntu)
Documentation=http://emby.media
After=network.target
Requires=network.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.embyu/bin/embyu.start
Restart=always

[Install]
WantedBy=multi-user.target

